# priv/container/release/supervisord_gate.conf
[supervisord]
nodaemon=true

[program:swm-gate]
command=/usr/local/bin/python3 /usr/local/lib/python/site-packages/swmcloudgate/start.py
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/swm-cloud-gate.log
